# crontab
00 02 * * 6	make -C /var/www/OPUS-API update-db-with-backup > /var/www/OPUS-API/update_db.log 2>&1

